# App/Cont1/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:uswgi]
command=/var/www/Cont1/venv/bin/uwsgi --ini /var/www/Cont1/nginx/uwsgi.ini
